# docker/logstash/etc/supervisord.d/logstash.ini
[program:logstash]
command=/usr/share/logstash/bin/logstash "--path.settings" "/usr/share/logstash/config/"
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/var/log/logstash.log
